# server-tools/roles/platform/templates/wireguard.service.j2
# Managed by Ansible (platform role) - do not edit on the host.
[Unit]
Description=WireGuard VPN (DHIS2 hosting platform)
After=docker.service traefik.service
Requires=docker.service

[Service]
Type=simple
WorkingDirectory={{ deploy_dir }}
ExecStart=/usr/bin/docker compose -p wireguard -f overlays/wireguard/docker-compose.yml up
ExecStop=/usr/bin/docker compose -p wireguard -f overlays/wireguard/docker-compose.yml down
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
